# package/usr/lib/systemd/user/odoo-companion.service
[Unit]
Description=Odoo Companion background polling service
# Order after the graphical session when one exists, but do NOT Wants= it: a
# Wants=graphical-session.target here gets pulled into the GDM greeter's own
# session and can break the login screen (black screen / no login prompt).
# After= alone is ordering-only and is ignored when the target isn't present.
After=graphical-session.target
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/usr/bin/odoo-companion-service
# on-failure (not always) so the clean "not a real user session" exit and the
# "autostart disabled" exit don't get restarted in a tight loop.
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
